# SingleOne_Backend/deploy/linux/singleone-api.service
[Unit]
Description=SingleOne API (.NET 6, sem Docker)
After=network.target postgresql.service

[Service]
WorkingDirectory=/opt/singleone-api-publish
ExecStart=/usr/bin/dotnet /opt/singleone-api-publish/SingleOneAPI.dll
Restart=always
RestartSec=10
SyslogIdentifier=singleone-api

# Ambiente
Environment=ASPNETCORE_ENVIRONMENT=Production

# URL do site (ajuste conforme seu domínio/IP)
# Se usar domínio: Environment=SITE_URL=https://seudominio.com
# Se usar IP: Environment=SITE_URL=http://SEU_IP_AQUI
Environment=SITE_URL=http://SEU_IP_AQUI

# As configurações de banco (DB_HOST, DB_USER, DB_PASSWORD, etc.)
# são carregadas do arquivo .env em /opt/singleone-api-publish
# via DotNetEnv (ver Startup.cs).

[Install]
WantedBy=multi-user.target



